# configs/home/.config/runit/sv/espanso/run
#!/bin/sh

USER=atalariq
HOME="/home/$USER"
[ -x "$HOME/Applications/bin/espanso-0.7.3" ] && "$HOME/Applications/bin/espanso-0.7.3" daemon
