# configs/linux/.config/runit/sv/espanso/run
#!/bin/sh

USER=atalariq
HOME=/home/$USER

ESPANSO_BIN=$HOME/Applications/bin/espanso-0.7.3
exec sh -c "$ESPANSO_BIN daemon && $ESPANSO_BIN restart"
