# configs/home/.config/runit/service/espanso/run
#!/bin/sh

# Start espanso daemon
espanso daemon
