# configs/home/.config/runit/service/espanso/run
#!/bin/sh

exec 2>&1
exec espanso daemon > /dev/null 2>&1
