# configs/linux/.config/runit/sv/custom-mappings/run
#!/bin/sh

USER=atalariq
HOME=/home/$USER

XMODMAP=$HOME/.Xmodmap

[ -f "$XMODMAP" ] && exec sh -c "xmodmap $HOME/.Xmodmap && xcape -e 'Mode_switch=Escape'" || exit 0
